# etc/supervisor.d/unjaild.ini
[program:unjaild]
command=/usr/local/bin/unjaild.sh
user=cosmovisor
autostart=false
autorestart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0